# scripts/systemd/retry-daemon.service
[Unit]
Description=Retry Metrics Collection Daemon
Documentation=https://github.com/user/retry
After=network.target
Wants=network.target

[Service]
Type=simple
User=retry
Group=retry
ExecStart=/usr/local/bin/retryd -config /etc/retry/daemon.json
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/retry /var/log/retry /tmp
CapabilityBoundingSet=
AmbientCapabilities=
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=retryd

# Working directory
WorkingDirectory=/var/lib/retry

# Environment
Environment=RETRY_LOG_LEVEL=info

[Install]
WantedBy=multi-user.target