# scripts/launchd/com.patience.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.patience.daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/patienced</string>
        <string>-config</string>
        <string>/usr/local/etc/patience/daemon.json</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/patience/daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/patience/daemon.error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/var/lib/patience</string>
    
    <key>UserName</key>
    <string>_patience</string>
    
    <key>GroupName</key>
    <string>_patience</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>RETRY_LOG_LEVEL</key>
        <string>info</string>
    </dict>
    
    <key>SoftResourceLimits</key>
    <dict>
        <key>NumberOfFiles</key>
        <integer>65536</integer>
        <key>NumberOfProcesses</key>
        <integer>4096</integer>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>5</integer>
    
    <key>ExitTimeOut</key>
    <integer>30</integer>
</dict>
</plist>